# mergeinputs-restart.path
[Path]
Unit=mergeinputs-restart.service
PathChanged=/dev/input/by-path

[Unit]
StartLimitIntervalSec=5
StartLimitBurst=20

[Install]
WantedBy=multi-user.target
