# mergeinputs-restart.path
[Path]
Unit=mergeinputs-restart.service
PathChanged=/dev/input/by-id

[Install]
WantedBy=multi-user.target
